# artifacts/docker/jupyter/rootfs/etc/supervisor/conf.d/jupyter.conf
[program:jupyter]
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
command=python3 -m jupyterlab --config=/etc/idekube/jupyter_server_config.py
user=idekube
directory=/home/idekube
